# udev/01-hdmi-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/osipov/.Xauthority", RUN+="/home/osipov/.udev/hdmi-hotplug"
